# config/systemd/user/dropbox.service
[Unit]
Description=Dropbox - file synchronization utility

[Service]
ExecStart=/usr/bin/dropbox
Restart=on-abort

[Install]
WantedBy=default.target
